# backend/scripts/lottery-api.service
[Unit]
Description=Scratchcard Website API Service
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/root/lottery-website/backend
Environment="PATH=/root/lottery-website/backend/.venv/bin"
ExecStart=/root/lottery-website/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
